# root/etc/s6-overlay/s6-rc.d/prisma/run
#!/command/with-contenv sh
cd /webui
npx prisma migrate deploy
npm run prisma:generate