# root/etc/s6-overlay/s6-rc.d/webui/run
#!/command/with-contenv sh
cd /
node webui